# docker/griffin-base/config/tomcat
#!/bin/bash

# description: Tomcat Start Stop Restart
# processname: tomcat
# chkconfig: 234 20 80
JAVA_HOME=/usr/java/latest
export JAVA_HOME
PATH=$JAVA_HOME/bin:$PATH
export PATH
JAVA_OPTS="-Dfile.encoding=UTF-8 -Xms128m -Xmx1024m -XX:PermSize=128m -XX:MaxPermSize=512m"
export JAVA_OPTS
echo $JAVA_OPTS
CATALINA_HOME=/apache/tomcat
RETVAL=$?
case "$1" in
start)
$CATALINA_HOME/bin/startup.sh
;;
restart)
$CATALINA_HOME/bin/shutdown.sh
$CATALINA_HOME/bin/startup.sh
;;
stop)
$CATALINA_HOME/bin/shutdown.sh
esac
exit 0
